# chaddesk/.config/systemd/user/change_wallpaper.timer
[Unit]
Description=Run change wallpapers hourly
After=xorg.target

[Timer]
OnBootSec=30s
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target
