# chaddesk/.config/systemd/user/change_wallpaper.timer
[Unit]
Description=Run change wallpapers hourly

[Timer]
OnUnitActiveSec=1h 

[Install]
WantedBy=timers.target
